# cmd/daemon/openlawsvpn-daemon.service
[Unit]
Description=openlawsvpn VPN daemon
Documentation=https://openlawsvpn.com
After=network.target dbus.service

[Service]
Type=simple
ExecStart=/usr/libexec/openlawsvpn-daemon
Restart=on-failure
RestartSec=3

# Allow TUN creation and routing without running as root.
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Hardening — drop everything else.
NoNewPrivileges=true
# PrivateUsers=no is explicit: with PrivateUsers=yes (systemd default for user services
# on some distros) the daemon runs in a child user namespace where CAP_NET_ADMIN is
# scoped to that namespace, not the host net namespace — TUNSETIFF fails with EPERM.
PrivateUsers=no
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ConfigurationDirectory=openlawsvpn
ReadWritePaths=%h/.config/openlawsvpn

[Install]
WantedBy=default.target
